# thirdparty/test-runner-image/inittab
::wait:-/startup
::ctrlaltdel:/bin/busybox reboot
::shutdown:/bin/busybox umount -a -r
::restart:/bin/busybox init
